# data/system/libvirt-dbus.service.in
#  SPDX-License-Identifier: LGPL-2.1+
#
#  This file is part of systemd.
#
#  systemd is free software; you can redistribute it and/or modify it
#  under the terms of the GNU Lesser General Public License as published by
#  the Free Software Foundation; either version 2.1 of the License, or
#  (at your option) any later version.

[Unit]
Description=Libvirt DBus Service

[Service]
BusName=org.libvirt
DynamicUser=yes
User=@SYSTEM_USER@
Group=@SYSTEM_USER@
SupplementaryGroups=libvirt
ExecStart=@sbindir@/libvirt-dbus --system

[Install]
Alias=org.libvirt.service
